# xdcterm.service.example
[Unit]
Description=xdcterm Delta Chat bot
After=network-online.target

[Service]
Type=simple
Environment=TERM=xterm-256color
WorkingDirectory=__WORKDIR__/backend
ExecStart=__UV_BIN__ run --directory __WORKDIR__/backend bot.py serve
Restart=on-failure
RestartSec=10
KillMode=mixed

[Install]
WantedBy=default.target
